# root/etc/crontabs/abc
0 3 * * * /app/dispatcher.sh >> /config/log/diskover/dispatcher.log 2>&1
