# root/etc/services.d/diskover_dispatcher/run
#!/usr/bin/with-contenv bash

cd /app/diskover || exit

exec \
	s6-setuidgid abc /usr/bin/python3 ./diskover.py -d /data -a
