# root/etc/services.d/socket_server/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/bin/python /app/diskover/diskover.py -l
